# deploy/guzel-hosting/supervisor-reverb.conf
; aaPanel → Supervisor → Add → veya /etc/supervisor/conf.d/human-reverb.conf
; Proje yolunu kendi sunucunuza göre düzenleyin (www kullanıcısı aaPanel varsayılanı)

[program:human-reverb]
process_name=%(program_name)s
command=php /www/wwwroot/humansocialpeople.com/artisan reverb:start --host=127.0.0.1 --port=8080
directory=/www/wwwroot/humansocialpeople.com
autostart=true
autorestart=true
stopwaitsecs=10
user=www
numprocs=1
redirect_stderr=true
stdout_logfile=/www/wwwroot/humansocialpeople.com/storage/logs/reverb.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=3
